# Configuration_Debian_Asus_Zenbook_UX430U/rc.local
#!/bin/sh -e

# Starting firewall
if [ -e '/etc/firewall.sh' ]
then
    /bin/sh '/etc/firewall.sh'
fi

exit 0
